# cron/cleanup.cron
59	23	*	*	*	/usr/local/bin/ftp_cleanup.pl
